# systemd/rebooter.timer
[Unit]
Description=Rebooter Timer

[Timer]
OnCalendar=Sat 5:00

[Install]
WantedBy=timers.target
