# systemd/rebooter.service
[Unit]
Description=Rebooter Service

[Service]
Type=oneshot
ExecStart=systemctl reboot
